# run_api.service
[Unit]
Description=Mongo DB

[Service]
ExecStart=/opt/run_api.sh

[Install]
WantedBy=multi-user.target
